# install/caddy-reload.service
[Unit]
Description=Reload Caddy

[Service]
Type=oneshot
ExecStart=/bin/systemctl try-reload-or-restart caddy.service