# install/caddy-reload.path
[Unit]
Description = Reload Caddy on CaddyFile Change

[Path]
PathChanged=/apps/TygerCaddy/TygerCaddy/data/caddyfile.conf

[Install]
WantedBy=multi-user.target