# install/gunicorn.socket
[Unit]
Description=gunicorn3 socket

[Socket]
ListenStream=/run/gunicorn3/socket

[Install]
WantedBy=sockets.target